# consul-agent/rootfs/etc/services.d/consul/run
#!/usr/bin/env bashio

CONSUL_SERVER_ADDR=$(bashio::config 'server_addr')
CONSUL_TOKEN=$(bashio::config 'consul_token')
CONSUL_BIND_INTERFACE=$(bashio::config 'bind_interface')
CONSUL_DC=$(bashio::config 'datacenter')
CONSUL_NAME=$(bashio::config 'consul_name')
LOG_LEVEL=$(bashio::config 'log_level')

SSL_ENABLED=$(bashio::config 'ssl')
CACERT_FILE=/ssl/$(bashio::config 'cafile')
CERT_FILE=/ssl/$(bashio::config 'certfile')
KEY_FILE=/ssl/$(bashio::config 'keyfile')


# When SSL_ENABLED is set add the extra options for SSL support
#
if [ $SSL_ENABLED == "true" ]; then
  if [ -e $CACERT_FILE -a -e $CERT_FILE -a -e $KEY_FILE ]; then
    continue
  fi

  bashio::exit.nok "When enabling SSL you should add the appropiate cafile, certfile, and keyfile to /ssl path..."
fi

CONSUL_BIND=
if [ ! -z "$CONSUL_BIND_INTERFACE" ]; then
  CONSUL_BIND_ADDRESS=$(ip -o -4 addr list $CONSUL_BIND_INTERFACE | head -n1 | awk '{print $4}' | cut -d/ -f1)
  if [ -z "$CONSUL_BIND_ADDRESS" ]; then
    bashio::exit.nok "Could not find IP for interface '$CONSUL_BIND_INTERFACE', exiting"
  fi

  CONSUL_BIND="-bind=$CONSUL_BIND_ADDRESS"
  bashio::log.info "==> Found address '$CONSUL_BIND_ADDRESS' for interface '$CONSUL_BIND_INTERFACE', setting bind option..."
fi

exec /consul agent \
  -config-dir=/etc/consul.d \
  -node=$CONSUL_NAME \
  -datacenter=$CONSUL_DC \
  -retry-join=$CONSUL_SERVER_ADDR \
  -log-level=$LOG_LEVEL \
  $CONSUL_BIND

